# aeroplot.cron
0 * * * * /usr/local/bin/python2.7 /Users/aayushchaturvedi/Projects/SandBox/Aeroplot/scraper.py >> ~/cron.log 2>&1